# etc/services.d/stdout/run
#!/usr/bin/execlineb -P
redirfd -u 0 /var/run/stdout cat
